# crontab
*/5 * * * * /home/jars99/mine/check.pl 2>&1 >> /var/log/check.log
