# rc.local
sleep 30s
/home/jars99/mine/startMining 2>&1 >> /var/log/mine/startMining
exit 0
